# etc/cron.d/apt-mirror
42 * * * *	run-one apt-mirror /etc/apt/mirrors.list
