# etc/cron.d/orange-box-mirror
42 * * * *	root	run-one apt-mirror /etc/apt/orange-box-mirror.list >>/var/log/orange-box-mirror.log 2>&1
