# etc/cron.d/orange-box-mirror
42 * * * *	run-one apt-mirror /etc/apt/orange-box-mirror.list
