# etc/init/orange-box-amt-monitor.conf
# orange-box-amt-monitor - monitor the AMT port

description	"Orange Box AMT Monitor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec orange-box-amt-monitor
